# debian/postinst
#!/bin/sh
set -e
case "$1" in
    configure)
        mkdir -p /usr/local/bin
        ln -sf /usr/bin/qrreader /usr/local/bin/qrreader
        echo "Installed /usr/bin/qrreader (symlinked at /usr/local/bin/qrreader)"
        echo "Per user: systemctl --user enable --now qrreader.service"
        ;;
esac
#DEBHELPER#
exit 0
