# debian/postinst
#!/bin/sh
set -e

case "$1" in
    configure)
        if [ -d /run/systemd/system ]; then
            echo "QR Reader installed to /usr/bin/qrreader"
            echo "Per-user service (run as each desktop user):"
            echo "  systemctl --user daemon-reload"
            echo "  systemctl --user enable --now qrreader.service"
            echo "  journalctl --user -u qrreader.service -f"
        fi
        ;;
esac

#DEBHELPER#

exit 0
